# roles/common/common-system/files/ttyS0.conf
# Serial file maintained by ansible

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [!2345]

respawn
exec /sbin/getty -8 38400 ttyS0

